# openssh-chroot/etc/s6-overlay/s6-rc.d/init-openssh-chroot-config/run
#!/usr/bin/with-contenv bash

set -e

echo Setting ChrootDirectory
echo "ChrootDirectory /config" >> /etc/ssh/sshd_config
